# supervisord.conf
[supervisord]
nodaemon=true

[program:top]
command=top
